# server/etc/crontab
@reboot /home/platform/_build/prod/rel/prod/bin/prod start